# chromeos-base/chromeos-bsp-surface-pro6/files/hardware-fix/99-surface-cover-attached-init.rules
ACTION=="add", SUBSYSTEM=="usb", ATTR{product}=="Surface Type Cover", RUN+="reload_i2c.sh"
